# profile.d/adm_bin.sh
../../../share/xworld/local/etc/profile.d/adm_bin.sh