# profile.d/bash-completion.sh
../../../share/xworld/local/etc/profile.d/bash-completion.sh