# profile.d/ccache_path.sh
../../../share/xworld/local/etc/profile.d/ccache_path.sh