# profile.d/editor.sh
../../../share/xworld/local/etc/profile.d/editor.sh