# profile.d/path_adm_bin.sh
../../../share/xworld/local/etc/profile.d/path_adm_bin.sh