# profile.d/path_xworld.sh
../../../share/xworld/local/etc/profile.d/path_xworld.sh