# profile.d/shell-goodies.sh
../../../share/xworld/local/etc/profile.d/shell-goodies.sh