# profile.d/upg_umask.sh
../../../share/xworld/local/etc/profile.d/upg_umask.sh