# profile.d/xworld.sh
../../../share/xworld/local/etc/profile.d/xworld.sh