# cron.d/decimate-logfiles.cron
#! /bin/sh
# Delete excess log files.

set -- --version 2 \
	--namespace decimate-logfiles_fl8bnlzry000fvzqtjg4bu8m0 \
	LOGSIZE_EBUILDS LOGSIZE_ELOGS
. /usr/local/bin/xworld/functions/lookup.sh
/usr/local/bin/xworld/decimate-files -f daemon /var/log/ebuilds "$LOGSIZE_EBUILDS"
/usr/local/bin/xworld/decimate-files -f daemon /var/log/ebuilds/elog "$LOGSIZE_ELOGS"
