# cron.d/decimate-binpkgs.cron
#! /bin/sh
# Delete excess cached distribution archives.

set -- --version 2 \
	--namespace decimate-binpkgs_t96tkd6t8t5gmmb86t293kayt \
	BINPKG_KEEP_SIZE
. /usr/local/bin/xworld/functions/lookup.sh
dir=/usr/portage/packages/All
test ! -d "$dir" && exit
/usr/local/bin/xworld/decimate-files -f daemon -ad 1 \
	"$dir" "$BINPKG_KEEP_SIZE"
