# cron.d/decimate-distcache.cron
#! /bin/sh
# Delete excess cached distribution archives.

set -- --version 2 DISTFILES_KEEP_SIZE
. /usr/local/bin/xworld/functions/lookup.sh
/usr/local/bin/xworld/decimate-files -f daemon -ad 1 \
	/usr/portage/distfiles "$DISTFILES_KEEP_SIZE"
