# overlay/etc/init.d/localswap
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License v2
# $Header: /etc/init.d/RCS/localswap,v 1.1 2005/08/04 22:42:46 root Exp $

depend() {
	need localmount
}

start() {
	ebegin "Activating shared Windows XP pagefile as swap device"
	mkswap /var/lib/swap/pagefile.sys > /dev/null
	swapon /var/lib/swap/pagefile.sys
	eend $? "Could not activate shared swap device"
}

stop() {
	ebegin "Deactivating swap in shared Windows XP pagefile"
	swapoff /var/lib/swap/pagefile.sys
	eend $? "Failed to deactivate shared swap device"
}
