# overlay/etc/init.d/vmware-server
#!/sbin/runscript

# Use this script and /etc/init.d/vmware-authd
# instead of the Gentoo-provided /etc/init.d/vmware script
# for this system!


# Source the distribution-provided script.
. /etc/init.d/vmware

# Override!
depend() {
	need vmware-authd
	use logger net
	after samba
}
