# cron.d/unused/decimate-binpkgs.cron
#!/bin/sh
# Delete excess cached distribution archives.
#
# $HeadURL: /caches/xsvn/uxadm/trunk/usr/local/etc/cron.d/daily/decimate-binpkgs.cron $
# $Author: root(xtreme) $
# $Date: 2006-08-23T03:52:22.370203Z $
# $Revision: 220 $

. /usr/local/bin/functions/lookup.sh --version 1 BINPKG_KEEP_SIZE
/usr/local/bin/decimate-files -f daemon -ad 1 \
	/usr/portage/packages/All "$BINPKG_KEEP_SIZE"
